# package/gluon-wan-dnsmasq/files/etc/hotplug.d/iface/50-gluon-wan-dnsmasq
if [ "$INTERFACE" = 'wan' -o "$INTERFACE" = 'wan6' ]; then
   /lib/graviton/wan-dnsmasq/update.lua
fi
