# package/graviton-admin-theme/files/etc/uci-defaults/luci-theme-graviton
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Graviton=/luci-static/graviton
        commit luci
EOF
	
