# package/graviton-core/files/etc/uci-defaults/zzz-gluon-upgrade
#!/bin/sh

for script in /lib/graviton/upgrade/*; do
	"$script"
done
